# DBManage/update.timer
[Unit]
Description= update BDD with weather fake data

[Timer]
OnUnitActiveSec=60

[Install]
WantedBy=timers.target